# roles/jeff/files/jeff-service.conf
description "jeff.conf"
author "Markus Klepp - git@mklepp.com"

start on started mountall
stop on shutdown

script
  exec /usr/local/bin/image-retriever.py
end script
